# devices/ZP120/filesystems/diagfs/etc/inittab
# autogenerated by gen_inittab.py for ARCH ppc; DO NOT EDIT
::sysinit:/etc/Configure
null::respawn:/wifi/netstartd
null::respawn:/etc/runanacapa
/dev/tts/0::askfirst:/bin/sh
::ctrlaltdel:/sbin/reboot
::shutdown:/etc/init.d/rcK
::restart:/sbin/init
